# deploy/com.bullbot.main.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bullbot.main</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>/Users/danield.runion/Bull-Bot/.venv/bin/python</string>
        <string>-m</string>
        <string>bullbot.main</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/danield.runion/Bull-Bot</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/danield.runion/Bull-Bot</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/danield.runion/Bull-Bot/logs/bullbot.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/danield.runion/Bull-Bot/logs/bullbot.stderr.log</string>
</dict>
</plist>
